# backend/src/scrape/scrape.dev.cron
* * * * * cd /app ; /usr/local/go/bin/go run /app/scrape/main.go > /app/scrape/debug.log 2>&1